# source/overlay/etc/init.d/esphosted
#!/sbin/openrc-run

# Copyright (c) 2025 Lukas Schmid

description="Switch ESP32 to Flash Boot mode"

depend()
{
	before modules
}

start() {
	yesno \$rc_verbose && verbose=yes

	ebegin "Starting ESP32 from Flash"
	eindent
	gpioset -t 0 ESP_nBOOT=1 # Set Bootmode Flash
	gpioset -t 500,0 ESP_nRST=0; # Reset ESP32
	eoutdent

	eend $?
}